# tor/services2/tor10.service
[Unit]
Description=Tor Daemon 10

[Service]
ExecStart=/usr/bin/tor -f /etc/tor/torrc.10

[Install]
WantedBy=multi-user.target