# tor/services2/tor13.service
[Unit]
Description=Tor Daemon 13

[Service]
ExecStart=/usr/bin/tor -f /etc/tor/torrc.13

[Install]
WantedBy=multi-user.target