# tor/services2/tor14.service
[Unit]
Description=Tor Daemon 14

[Service]
ExecStart=/usr/bin/tor -f /etc/tor/torrc.14

[Install]
WantedBy=multi-user.target