# tor/services2/tor4.service
[Unit]
Description=Tor Daemon 4

[Service]
ExecStart=/usr/bin/tor -f /etc/tor/torrc.4

[Install]
WantedBy=multi-user.target