# tor/services2/tor7.service
[Unit]
Description=Tor Daemon 7

[Service]
ExecStart=/usr/bin/tor -f /etc/tor/torrc.7

[Install]
WantedBy=multi-user.target