# tor/services/tor2.service
[Unit]
Description=Tor Daemon 2

[Service]
ExecStart=/usr/bin/tor -f /root/tor/torrc.2

[Install]
WantedBy=multi-user.target