# tor/services/tor7.service
[Unit]
Description=Tor Daemon 7

[Service]
ExecStart=/usr/bin/tor -f /root/tor/torrc.7

[Install]
WantedBy=multi-user.target